# deploy/systemd/shopvivaliz-token-renewer.service
[Unit]
Description=ShopVivaliz Olist OAuth Token Renewer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/site-shopvivaliz
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/ubuntu/site-shopvivaliz/daemon-token-renewer.py --interval 7200 --retry-interval 900
Restart=always
RestartSec=60
TimeoutStopSec=30
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
