# deploy/systemd/shopvivaliz-sync-products.service
[Unit]
Description=ShopVivaliz Olist Catalog Detail Sync
After=network-online.target shopvivaliz-token-renewer.service
Wants=network-online.target shopvivaliz-token-renewer.service

[Service]
Type=simple
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/site-shopvivaliz
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/ubuntu/site-shopvivaliz/daemon-sync-products.py --interval 300 --workers 4
Restart=always
RestartSec=30
TimeoutStopSec=30
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
